# deploy_tools/superlists.conf
description "superlists start"

start on net-device-up
stop on shutdown

respawn

setuid ubuntu
chdir /home/ubuntu/tutorial-clean-python

exec gunicorn --bind unix:/tmp/superlists.socket superlists.wsgi:applicaiton
